# examples/mvad-up@.service
[Unit]
Description=mvad tunnel to %i
After=network-online.target mvad-lockdown.service
Wants=network-online.target
StartLimitIntervalSec=30
StartLimitBurst=5

[Service]
Type=simple
ExecStart=/usr/bin/mvad up %i
Restart=on-failure
RestartSec=5s
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
DefaultInstance=se-got-wg-001
