# examples/mvad-check.service
[Unit]
Description=probe the mvad tunnel and fail over

[Service]
Type=oneshot
# KillMode: reconnect may spawn transport shims that must outlive this unit.
KillMode=process
ExecStart=/usr/bin/mvad reconnect --if-dead
