# examples/mvad-lockdown.service
[Unit]
Description=mvad lockdown firewall
Before=network-pre.target
Wants=network-pre.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/nft -f /var/lib/mvad/lockdown.nft
ExecStop=/usr/bin/true
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
