# examples/mvad-up@.service
[Unit]
Description=mvad tunnel to %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/mvad up %i
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
